# supervisord.conf
[program:manage]
command: /usr/bin/python -m manage
directory: /home/ubuntu/cfpydev-imagr
autostart: true