# supervisord.conf
[program:manage]
command: /usr/bin/python manage.py runserver
directory: /home/ubuntu/cfpydev-imagr
autostart: true
stdout_logfile: /tmp/djangolog
stderr_redirect: true